# infra/systemd/canton-tunnel.service.example
# Canton + Splice SSH tunnel — run on VPS 2 only (not on validator VPS 1).
#
# Install:
#   sudo cp infra/systemd/canton-tunnel.service.example /etc/systemd/system/canton-tunnel.service
#   sudo nano /etc/systemd/system/canton-tunnel.service   # verify Docker IPs on VPS 1
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now canton-tunnel
#
# Re-check IPs on validator (162.250.190.204 TestNet):
#   docker inspect -f '{{range .NetworkSettings.Networks}}{{.IPAddress}}{{end}}' splice-validator-participant-1
#   docker inspect -f '{{range .NetworkSettings.Networks}}{{.IPAddress}}{{end}}' splice-validator-nginx-1

[Unit]
Description=CanQuest SSH tunnel to Canton participant + Splice nginx (TestNet)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
# PARTICIPANT_IP=172.18.0.5  NGINX_IP=172.18.0.7  (TestNet May 2026 — verify with docker inspect)
ExecStart=/usr/bin/ssh -N \
  -o ServerAliveInterval=30 \
  -o ServerAliveCountMax=3 \
  -o ExitOnForwardFailure=yes \
  -o StrictHostKeyChecking=accept-new \
  -L 127.0.0.1:7575:172.18.0.5:7575 \
  -L 127.0.0.1:8080:172.18.0.7:80 \
  root@162.250.190.204
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
