# infra/systemd/canton-tunnel.service
# Persistent SSH tunnel from VPS 2 → VPS 1 (Canton validator)
#
# VPS 1 nginx config (confirmed 2025-05):
#   server_name wallet.localhost        → validator:5003  (Splice API)
#   server_name json-ledger-api.localhost → participant:7575 (Ledger API)
#   nginx binds to 127.0.0.1:80 on VPS 1 host
#
# Tunnels:
#   VPS2:7575 → VPS1 participant container 172.19.0.5:7575  (Canton JSON Ledger API langsung)
#   VPS2:8080 → VPS1 127.0.0.1:80                          (Splice nginx — pakai Host header)
#
# Cara akses setelah tunnel aktif (dari VPS 2):
#   curl http://127.0.0.1:7575/livez                                          → Ledger API
#   curl -H "Host: wallet.localhost" http://127.0.0.1:8080/api/validator/v0/admin/users → Splice API
#
# SETUP on VPS 2:
# 1. Generate SSH key: ssh-keygen -t ed25519 -f /root/.ssh/canton_tunnel -N ""
# 2. Tambah public key ke VPS 1: echo "$(cat /root/.ssh/canton_tunnel.pub)" >> /root/.ssh/authorized_keys
# 3. Copy file ini ke: /etc/systemd/system/canton-tunnel.service
# 4. Enable: systemctl daemon-reload && systemctl enable --now canton-tunnel

[Unit]
Description=SSH tunnel to Canton Validator VPS 1 (162.250.191.46)
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Restart=always
RestartSec=10

ExecStart=/usr/bin/ssh \
  -N \
  -o ServerAliveInterval=30 \
  -o ServerAliveCountMax=3 \
  -o ExitOnForwardFailure=yes \
  -o StrictHostKeyChecking=accept-new \
  -o BatchMode=yes \
  -i /root/.ssh/canton_tunnel \
  -L 7575:172.19.0.5:7575 \
  -L 8080:127.0.0.1:80 \
  root@162.250.191.46

[Install]
WantedBy=multi-user.target
